# crontab
*	*	*	*	*	sh Users/tejbirwason/Desktop/githubauto.sh
